# launchd/com.devmini.orchestrator.dashboard-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.devmini.orchestrator.dashboard-server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Volumes/devssd/orchestrator/bin/dashboard_server.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Volumes/devssd/orchestrator</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/devmini/Library/Logs/devmini/dashboard-server.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/devmini/Library/Logs/devmini/dashboard-server.stderr.log</string>
</dict>
</plist>
